# sv/master/run
#!/bin/sh

exec buildbot start --nodaemon ../../master
