# sv/worker/run
#!/bin/sh

exec buildbot-worker start --nodaemon ../../worker1
